# root/etc/uci-defaults/99-acmesh-console-cleanup
#!/bin/sh

rm -f /www/luci-static/resources/view/acmesh/operations.js
mkdir -p /etc/acmesh-console /etc/acmesh-console/ssh
chmod 700 /etc/acmesh-console /etc/acmesh-console/ssh
exit 0
